# template-files/centos/7.0/etc/systemd/system/gq-get-userdata.service
[Unit]
Description=GQ Get Userdata Script
After=networking.target

[Service]
Type=simple
ExecStart=/usr/local/qstack/bin/gq-get-userdata